# files/etc/sv/runsvdir-pbui/run
#!/bin/sh
exec sudo -H -u pbui runsvdir -P /home/pbui/.service 'log:.............................................................................'
